# docker/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/supervisord.pid

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
autorestart = true
user = root